# rootfs/etc/services.d/app/run
#!/command/with-contenv bash
set -euo pipefail

ENV_FILE="/config/aio/generated.env"
if [ -f "${ENV_FILE}" ]; then
    set -a
    . "${ENV_FILE}"
    set +a
fi

exec python3 /usr/local/bin/aio-template-app.py
